# dotfiles/otelcol-contrib.service
[Unit]
Description=OpenTelemetry Collector (contrib)
After=network.target

[Service]
ExecStart=/usr/local/bin/otelcol-contrib --config /etc/otelcol-contrib/config.yaml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
